# datastore-hotplug-handler@.service
[Service]
Type=oneshot
ExecStart=/usr/local/sbin/pbs_datastore_hotplug_handler.sh %i
